# crontab
40 5 * * * rsync -avzu --delete /mnt/data1/encoded/ /mnt/data5/encoded/
*/15 * * * * /home/yude/switchbot-discord/post.sh